# debian/opencpu-server-files/cron.d/opencpu-apps
#This cronjob updates the user apps library every 5 minutes
*/5 * * * * root /usr/lib/opencpu/scripts/userapps.sh